# support/aws-kinesis-agent.cron
* * * * * root service aws-kinesis-agent status || service aws-kinesis-agent start > /dev/null 2>&1
