# support/aws-kinesis-agent.cron
* * * * * root /usr/bin/aws-kinesis-agent-babysit > /dev/null 2>&1
